# airootfs/etc/systemd/system/localegen.service
[Unit]
Description=locale-gen

[Service]
Type=oneshot
ExecStart=/usr/bin/locale-gen

[Install]
WantedBy=multi-user.target